# systemd/semantic-backup.timer
[Unit]
Description=Run Semantic Backup 1 hour after boot

[Timer]
OnBootSec=1h
Persistent=true

[Install]
WantedBy=default.target
